# roles/mn-pay-faucet/templates/pay-faucet.service.j2
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/pay-faucet.sh

[Install]
WantedBy=default.target